# sync-scheduler/crontab
* * * * * root bash sync-call.sh
# Single container cron example